# roles/configure_supervisor/templates/nud.conf.j2
[program:nud]
command = nud -datadir=/home/nud/.nu
user = nud
redirect_stderr = true
stdout_logfile = /home/alp/logs/nud.log